# deploy/supervisor/stock-analyzer.conf
[program:stock-analyzer-cron]
command=/bin/bash -c "while true; do python /path/to/stock-daily-analyzer/main.py; sleep 86400; done"
directory=/path/to/stock-daily-analyzer
user=USER_TO_REPLACE
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/stock-analyzer/supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PATH="/path/to/stock-daily-analyzer/venv/bin:%(ENV_PATH)s"
